# Poke-Bowl---updated-January/deployment/chromium-kiosk.service
[Unit]
Description=Chromium Kiosk Mode for Poke Bowl Inventory
After=pokebowl-inventory.service graphical.target
Wants=graphical.target
Requires=pokebowl-inventory.service

[Service]
Type=simple
User=YOUR_USERNAME
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/YOUR_USERNAME/.Xauthority"

# Wait for web server to be ready
ExecStartPre=/bin/sleep 5

# Launch Chromium in kiosk mode
ExecStart=/usr/bin/chromium-browser \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-features=TranslateUI \
    --no-first-run \
    --disable-suggestions-service \
    --disable-save-password-bubble \
    --start-fullscreen \
    http://localhost:8080

Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target

