# Poke-Bowl---updated-January/deployment/pokebowl-inventory.service
[Unit]
Description=Poke Bowl Inventory System
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/Poke-Bowl---updated-January
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/YOUR_USERNAME/.Xauthority"

# Restart policy for production reliability
Restart=always
RestartSec=10

# Start the inventory system
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/Poke-Bowl---updated-January/backend/main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pokebowl-inventory

# Resource limits (adjust based on requirements)
MemoryMax=2G
CPUQuota=400%

[Install]
WantedBy=multi-user.target

